# ami-builder/scripts/ecp-boot.service
[Unit]
Description=Express Compute Cluster Bootstrap
After=network-online.target cloud-final.service
Wants=network-online.target
ConditionPathExists=!/opt/eks-d/.installation_complete

[Service]
Type=oneshot
RemainAfterExit=true
Environment=HOME=/root
WorkingDirectory=/opt/cluster-setup
ExecStart=/bin/bash /opt/cluster-setup/workstation-boot.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
